# eshop/db/mycat/program.conf
[supervisord]
nodaemon=true

[program:xinetd]
command=/etc/init.d/xinetd restart

[program:mycat]
command=/usr/local/mycat/bin/mycat console